# services/dynamodb-local.service
[Unit]
Description=DynamoDB Local for URL Shortener
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/urlshortener
Restart=always
RestartSec=10

# Start DynamoDB Local container
ExecStartPre=-/usr/bin/docker stop dynamodb-local
ExecStartPre=-/usr/bin/docker rm dynamodb-local
ExecStart=/usr/bin/docker run --rm \
  --name dynamodb-local \
  -p 8000:8000 \
  -v /var/lib/dynamodb-local:/data \
  amazon/dynamodb-local:latest \
  -jar DynamoDBLocal.jar -sharedDb -dbPath /data

# Stop container
ExecStop=/usr/bin/docker stop dynamodb-local

# Restart policy
TimeoutStartSec=0
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
