# frontend/deployment/urlshortener-frontend.service
[Unit]
Description=URL Shortener Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/urlshortener
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="NEXT_PUBLIC_API_GATEWAY=https://api.masondrake.dev/10012"
Environment="NEXT_PUBLIC_QRCODE_FUNCTION=https://api.masondrake.dev/10012/function/qrcode-go"
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
